# deploy/eu5modplanner.service
[Unit]
Description=EU5 Mod Planner
After=network.target docker.service
Requires=docker.service

[Service]
WorkingDirectory=__DEPLOY_PATH__
ExecStart=/usr/bin/dotnet __DEPLOY_PATH__/Eu5ModPlanner.dll
Restart=always
RestartSec=5
KillSignal=SIGINT
SyslogIdentifier=eu5modplanner
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=ASPNETCORE_URLS=http://0.0.0.0:5068
Environment=ConnectionStrings__PlannerDatabase=Host=127.0.0.1;Port=54329;Database=eu5_mod_planner;Username=eu5planner;Password=CHANGE_ME

[Install]
WantedBy=multi-user.target
